# docker/jobs/crontab
* * * * * root php bin/console crypto:wallet:update
