# docker/jobs/crontab
*/5 * * * * php /application/bin/console crypto:wallet:update
